# warehouse/systemd/donorpipe-ondemand.service
[Unit]
Description=DonorPipe on-demand refresh poller
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=<pi-user>
WorkingDirectory=/home/<pi-user>/donorpipe_warehouse/donorpipe
EnvironmentFile=-/home/<pi-user>/donorpipe_warehouse/donorpipe/.env
Environment="PATH=/home/<pi-user>/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/<pi-user>/donorpipe_warehouse/donorpipe/warehouse/poll_refresh.sh \
  --config warehouse/warehouse_pi_config.json
Restart=on-failure
RestartSec=30
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
