# warehouse/systemd/donorpipe-nightly.service
[Unit]
Description=DonorPipe nightly data update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=<pi-user>
WorkingDirectory=/home/<pi-user>/donorpipe_warehouse/donorpipe
EnvironmentFile=-/home/<pi-user>/donorpipe_warehouse/donorpipe/.env
Environment="PATH=/home/<pi-user>/.local/bin:/usr/local/bin:/usr/bin:/bin"
ExecStart=/home/<pi-user>/donorpipe_warehouse/donorpipe/warehouse/update.sh \
  --config warehouse/warehouse_pi_config.json nightly
StandardOutput=journal
StandardError=journal
