# warehouse/systemd/donorpipe-nightly.timer
[Unit]
Description=DonorPipe nightly data update timer

[Timer]
OnCalendar=*-*-* 03:00:00
RandomizedDelaySec=5min
Persistent=true

[Install]
WantedBy=timers.target
